# .ci/nextbsd/org.gershwin.ci-build.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- CI-only: builds gershwin-developer under the NextBSD kernel at boot, then
     powers off. Injected into the disposable QEMU image by prep-image.sh; it is
     never part of a shipping system. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.gershwin.ci-build</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>/build/ci-build.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/dev/console</string>
    <key>StandardErrorPath</key>
    <string>/dev/console</string>
</dict>
</plist>
